# etc/systemd/user/shell@.service
[Unit]
Requisite=shell.service
After=shell.service

[Service]
Type=oneshot
ExecStart=/usr/bin/shell %I
